# freebsd/rc.d/beadsd
#!/bin/sh
#
# PROVIDE: beadsd
# REQUIRE: NETWORKING
# KEYWORD: shutdown
#
# Single-writer beads service. ONE instance per project. rc.subr namespaces
# every variable by the script's filename, so install this script once per
# project under a distinct name (beadsd_mu, beadsd_at, …) — e.g. by symlink:
#
#   ln -s /usr/local/etc/rc.d/beadsd /usr/local/etc/rc.d/beadsd_mu
#   ln -s /usr/local/etc/rc.d/beadsd /usr/local/etc/rc.d/beadsd_at
#
# All per-instance settings (db, listen ip:port, repo) live in a TOML config:
# by default ~<user>/.config/beadsd/<instance>.toml, where <instance> is the
# script name minus the "beadsd_" prefix (beadsd_mu -> mu.toml). See
# config/*.example.toml in the beadsd repo.
#
# /etc/rc.conf (mu shown):
#   beadsd_mu_enable="YES"
#   beadsd_mu_user="tcovert"
#   # optional: beadsd_mu_config="/home/tcovert/.config/beadsd/mu.toml"
#   # optional: beadsd_mu_command, beadsd_mu_log, beadsd_mu_rust_log

. /etc/rc.subr

name=$(basename "$0")
rcvar="${name}_enable"

load_rc_config "$name"

eval ": \${${name}_enable:=\"NO\"}"
eval "_user=\"\${${name}_user:-}\""
eval "_home=\"\${${name}_home:-/home/${_user}}\""
eval "_command=\"\${${name}_command:-${_home}/.local/bin/beadsd}\""
eval "_log=\"\${${name}_log:-/var/log/${name}.log}\""
eval "_rust_log=\"\${${name}_rust_log:-beadsd=info}\""

# Per-instance config path: default ~<user>/.config/beadsd/<instance>.toml.
_instance=${name#beadsd_}
eval "_config=\"\${${name}_config:-${_home}/.config/beadsd/${_instance}.toml}\""

[ -n "${_user}" ] || err 1 "${name}_user is not set in /etc/rc.conf"
[ -f "${_config}" ] || err 1 "config not found: ${_config} (see config/*.example.toml)"

# rc.subr wraps the command with `su -m ${name}_user`, so daemon(8) runs as
# that user. We deliberately do NOT pass `daemon -u` (that + `-r` loops
# setusercontext() and crashes on restart — same trap as claude_proxy). HOME
# must be set explicitly: env(1) does not carry it through.
pidfile="/var/run/${name}/${name}.pid"

# PATH must include ~/.local/bin so beadsd finds br (and friends): daemon(8)
# does not inherit a login shell's PATH.
command="/usr/sbin/daemon"
command_args="-r -P ${pidfile} -o ${_log} \
    /usr/bin/env HOME=${_home} \
    PATH=${_home}/.local/bin:/usr/local/bin:/usr/bin:/bin \
    RUST_LOG=${_rust_log} \
    ${_command} --config ${_config}"

start_precmd="beadsd_prestart"
beadsd_prestart() {
    # Pidfile + log under a user-owned dir: the daemon runs as ${_user} (rc.subr
    # su -m) and cannot write /var/run (root:wheel 755).
    install -d -o "${_user}" -g "${_user}" -m 0755 "/var/run/${name}"
    install -o "${_user}" -g "${_user}" -m 644 /dev/null "${_log}" 2>/dev/null || \
        chown "${_user}:${_user}" "${_log}"
}

run_rc_command "$1"
